# com.pai-mobile.plist.example
<!-- PAI Mobile Gateway — launchd auto-start config

     To use:
     1. Copy this file: cp com.pai-mobile.plist.example ~/Library/LaunchAgents/com.pai-mobile.plist
     2. Replace __HOME__ with your home directory (e.g., /Users/yourname)
     3. Replace __PROJECT_DIR__ with the full path to this project
     4. Load: launchctl load ~/Library/LaunchAgents/com.pai-mobile.plist
-->
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pai-mobile</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>__PROJECT_DIR__/dist/index.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__PROJECT_DIR__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>__HOME__</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>StandardOutPath</key>
  <string>__PROJECT_DIR__/data/logs/launchd-stdout.log</string>

  <key>StandardErrorPath</key>
  <string>__PROJECT_DIR__/data/logs/launchd-stderr.log</string>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
